# deploy/systemd/tg-ws-proxy-update.service
[Unit]
Description=Update TG WS Proxy from git
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/tg-ws-vps/deploy/systemd/update-headless.sh
