# packaging/com.tg-ws-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tg-ws-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>tg-ws-proxy</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>1080</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/tg-ws-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tg-ws-proxy.log</string>
</dict>
</plist>